# debian/mobian-tweaks-common.mobian-usb-gadget.service
[Unit]
Description=Configure USB gadget

[Service]
Type=oneshot
ExecStart=/usr/sbin/mobian-usb-gadget setup
ExecStop=/usr/sbin/mobian-usb-gadget reset
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
